# launchd/com.briefing.recap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.briefing.recap</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/osascript</string>
		<string>-e</string>
		<string>do shell script "/usr/bin/python3 /Users/jeffreystclaire/Claude/morning-briefing/morning_briefing.py recap &gt; /tmp/briefing-recap.log 2&gt; /tmp/briefing-recap.err"</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/tmp/briefing-recap-osa.err</string>
	<key>StandardOutPath</key>
	<string>/tmp/briefing-recap-osa.log</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Weekday</key><integer>1</integer>
			<key>Hour</key><integer>14</integer>
			<key>Minute</key><integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key><integer>2</integer>
			<key>Hour</key><integer>14</integer>
			<key>Minute</key><integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key><integer>3</integer>
			<key>Hour</key><integer>14</integer>
			<key>Minute</key><integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key><integer>4</integer>
			<key>Hour</key><integer>14</integer>
			<key>Minute</key><integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key><integer>5</integer>
			<key>Hour</key><integer>14</integer>
			<key>Minute</key><integer>0</integer>
		</dict>
	</array>
</dict>
</plist>
